# deploy/aiengjobs-refresh.timer
[Unit]
Description=Run aiengjobs refresh nightly

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
